# systemd/mtr-server.service
[Unit]
Description=MTR Server
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/mtr-server -config /etc/mtr/server.yaml
Restart=always
RestartSec=5
User=mtr
Group=mtr
WorkingDirectory=/var/lib/mtr
StateDirectory=mtr
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
